# systemd/system/ydotoold.service
[Unit]
Description=ydotool daemon

[Service]
Type=simple
ExecStart=/usr/bin/ydotoold

[Install]
WantedBy=default.target
